# CelerisAgent/deploy/systemd/celeris-agent.service
[Unit]
Description=CelerisAgent API service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=celeris
Group=celeris
WorkingDirectory=/srv/celeris/current/CelerisAgent
EnvironmentFile=/etc/celeris-agent.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/srv/celeris/.venv/bin/python /srv/celeris/current/CelerisAgent/app.py --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
